# esphome-lightsd.service
# Reference service file — the installer (install.sh) generates the real
# ~/.config/systemd/user/esphome-lightsd.service with correct paths.
#
# Key differences from a system-level service:
#   - Installed as a USER service (systemctl --user ...), not system-wide.
#   - ExecStart uses the Python 3.11 venv created by the installer.
#   - Socket path is derived from XDG_RUNTIME_DIR automatically by the daemon.
#   - No EnvironmentFile= -- env loading is handled by esphome-lightsd.py.

[Unit]
Description=ESPHome Lights Daemon
Documentation=https://github.com/platima/ESPHome-Python
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/lib/esphome-lights/venv/bin/python %h/.local/lib/esphome-lights/esphome-lightsd.py
WorkingDirectory=%h/.local/lib/esphome-lights
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
